# liverootfs/etc/cron.d/czo-motd-czolsb
# Filename: czo-email-at-reboot

@reboot      root   /etc/czolsb > /etc/motd 2> /dev/null
